# ai-control-plane/scripts/oneflow-restore-drill.timer
[Unit]
Description=Weekly OneFlow restore drill (Sunday 03:30 UTC, after Saturday backup)
Documentation=file:///etc/systemd/system/oneflow-restore-drill.service

[Timer]
OnCalendar=Sun *-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=true
Unit=oneflow-restore-drill.service

[Install]
WantedBy=timers.target
